# bist-scheduler.service
[Unit]
Description=BIST AI System Automated Data Pipeline Scheduler
After=network.target postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bist-pattern/BIST-Ubuntu
ExecStart=/opt/bist-pattern/venv/bin/python /opt/bist-pattern/BIST-Ubuntu/scheduler_daemon.py
Environment="TZ=Europe/Istanbul"
Environment="PYTHONPATH=/opt/bist-pattern:/opt/bist-pattern/BIST-Ubuntu"

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security (optional)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
